# scripts/cookbooks/marine-planner/templates/default/app.conf.erb

description "gunicorn app"
author      "Edwin Knuth"


start on filesystem
stop on runlevel S
<% if node[:user] == "vagrant" %>
script
    export HOME="/vagrant/mp"
    export VENV="/usr/local/venv/marine-planner"
    cd $HOME
    exec sudo -u vagrant $VENV/bin/gunicorn -w 3 -b 0.0.0.0:8888 vagrant-wsgi
end script
<% else %>
script
    export HOME="/usr/local/apps/marine-planner/mp"
    export VENV="/usr/local/venv/marine-planner"
    cd $HOME
    exec sudo -u www-data $VENV/bin/gunicorn -w 3 -b 0.0.0.0:8888 wsgi
end script
<% end %>